# backend/assets/spaceship.service
[Unit]
Description=Spaceship Homepage

[Service]
ExecStart=/usr/bin/spaceship

[Install]
WantedBy=default.target
